# roles/hockeybot-service/files/lita-hockeybot.conf
description "lita-hockeybot"

start on filesystem or runlevel [2345]
stop on runlevel [!2345]


setuid BOT_USER
env HOME=/srv/repos/lita-slack-jira-hockeyapp-mode/workspace
env SLACK_TOKEN=

# for HockeyApp crashes
env HOCKEYAPP_TOKEN=
env HOCKEYAPP_URL="https://rink.hockeyapp.net/api/2/apps"
env HOCKEYAPP_TOPN_CRASHES=10
env HOCKEYAPPS_AND_VERSIONS="{'2980b0f3ffcf39cc0b562526d34bac95'=>221, '346d8093aa974667a7e615dbaca2fab9'=>95, '41391935d662401996bdeea9b7199764'=>51, 'bbd265f4aa695cda8e48b1b7d633bbf9'=>212 }"
env HOCKEYAPP_MASTER_APP_IDS='{"OD": "5caab2ce2b3df230f30ad8428e0329d3", "COZMO": "3e78b35bae854fa8b3137db388bed136", "COZMO_IOS": "714e67f2ff5248469ee4d8419d0dfd74", "OD_IOS": "99304009eb0a0fe6a34939617e4daf4b"}'
env PACKAGE_APPS='{"COZMO": "com.anki.cozmo", "OD": "com.anki.overdrive"}'

# for Release Schedule
env RELEASE_PAGE_ID="148477916"

# for JIRA tickets
env JIRA_SITE=https://ankiinc.atlassian.net
env JIRA_TIMEOUT=150
env JIRA_USERNAME=
env JIRA_PASSWORD=
env JIRA_COMPONENTS="crashes,errors,release checklist"
env STR_EXCLUDES="fault addr ,libmono.,Unknown.,BLECentralMultiplexer centralManager:didDiscoverPeripheral:advertisementData:RSSI,storage/emulated/0/Android/data/,line ,OutOfMemoryError: Failed to allocate a "
env JIRA_PROJECTS="{'od'=>'OD','overdrive'=>'OD','OverDrive'=>'OD','cozmo'=>'COZMO','Cozmo'=>'COZMO','CozmoOne'=>'COZMO','fox'=>'OD','Vic'=>'VIC','VIC'=>'VIC','vic'=>'VIC'}"
env JIRA_PROJECT_OPEN_TRANSITIONS="{'OD'=>'101','COZMO'=>'341',}"
#env JIRA_PROJECT_CLOSE_TRANSITIONS="{'OD'=>'331','COZMO'=>'331',}"
env JIRA_PROJECT_CLOSE_TRANSITIONS="{'OD'=>{'Open'=>['61'],'In Progress'=>['91','61'],'To Do'=>['61'],'Code Review'=>['121'],'Test Ready'=>['81','61']},'COZMO'=>{'Open'=>['331'],'In Progress'=>['91','331'],'Code Review'=>['121'],'Test Ready'=>['81','331']}}"
env CLOSE_TICKET_TIMEOUT="7"
env JIRA_MAX_RESULTS="1000"
env VIP_CATEGORY_COLUMN="{'name'=>'customfield_11404','value'=>{'KS'=>['101 (KS)','001 (KS)']}}"
env CRASHES_FILTER="{'cozmo-crashes'=>20101,'cozmo-last-24'=>30500,'cozmo-last-24-crashes'=>30203,'cozmo-last-24-failed'=>33002,'od-crashes'=>31400,'od-last-24'=>30202,'od-last-24-failed'=>33001,'cozmo-last-24-new-and-failed'=>38548,'od-last-24-new-and-failed'=>38549,'vic-last-24-new-and-failed'=>38616,'vic-last-24-new-and-failed-high-high'=>39118,'victor-blocked-lkg'=>38830,'victor-new-last-15'=>38705,'victor-new-last-60'=>38706,'victor-updated-last-15'=>38707,'victor-updated-last-60'=>38708}"
env EZ_SITE=https://anki.ezofficeinventory.com/assets
env EZ_TOKEN=

# for cherry_pick functionality
env GIT_URI_OD=git@github.com:/anki/overdrive-one.git
env GIT_URI_COZMO=git@github.com:/anki/cozmo-one.git
env GIT_TOKEN=
env GIT_API_URI_VECTOR="https://api.github.com/repos/anki/victor"
env TEAMCITY_URI=https://build.ankicore.com
env TEAMCITY_INFO="{'OD'=>['OverDriveOne_Dev_PullRequest','OverDriveOne_Dev_Android','OverDriveOne'], 'COZMO'=>['CozmoOne_Dev_PullRequestsIOS','CozmoOne_Dev_PullRequestAndroid','CozmoOne'], 'VIC'=>['','Victor_PullRequestAndroid','Victor']}"

# for Mode ERROR reports
env MODE_TOKEN=""
env MODE_PASSWORD=""
env MODE_TIME_WAIT="5400"
env MODE_REPORT_TIMEOUT="3500"
env QUERY_TIMEOUT="30"
env EMAIL_USERNAME="do-not-reply"
env EMAIL_PASSWORD="@nk1TEST"
env LABELS="errors"
env LOG_WARNING="NO"
env OD_PROD_OCCURENCE_THRESHOLD="15000"
env OD_DEV_OCCURENCE_THRESHOLD="100"
env COZMO_PROD_OCCURENCE_THRESHOLD="1000"
env COZMO_DEV_OCCURENCE_THRESHOLD="25"
env OD_DEV_MODE_URL="https://modeanalytics.com/anki/reports/deb8d35a9ac1?param_period=PERIOD&run=now"
env OD_PROD_MODE_URL="https://modeanalytics.com/anki/reports/30dd44e348c0?param_period=PERIOD&run=now"
env COZMO_DEV_MODE_URL="https://modeanalytics.com/anki/reports/26ede1626769?param_period=PERIOD&run=now"
env COZMO_PROD_MODE_URL="https://modeanalytics.com/anki/reports/170693d790ba?param_period=PERIOD&run=now"
env COZMO_DEV_MODE_URL_SPECIFIC_ERROR="https://modeanalytics.com/anki/reports/7c043203af6c"
env OD_DEV_MODE_URL_SPECIFIC_ERROR="https://modeanalytics.com/anki/reports/7c043203af6c"
env TO_EMAIL="ww-test@anki.com"
env CC_EMAIL="stephen@anki.com,jliptak@anki.com,paul@anki.com,ww-test@anki.com"

# for app info on Google Play and Apple Store
env URL_GOOGLEPLAY_API="https://play.google.com/store/apps/details?hl=en&id="
#env URL_GOOGLEPLAY_API="https://data.42matters.com/api/v2.0/android/apps/lookup.json?p="
env URL_APPLESTORE_API="http://itunes.apple.com/lookup?bundleId="
env ACCESS_GOOGLEPLAY_TOKEN=""

# for self updating
env UPDATE_SCRIPT="/srv/repos/ansible-bot-deployment/self_update.sh"
env ALLOWED_UPDATE_USERS="U17B8V7CZ,U040HGVJJ,U76PP1Q1Z,U4CBWEA2U"

# For get DAS version from MODE
env DEV_DAS_VERSION_MODE_URL="https://modeanalytics.com/anki/reports/88d38b9bf27c"
env PROD_DAS_VERSION_MODE_URL="https://modeanalytics.com/anki/reports/007c2ebacf18"
env BETA_DAS_VERSION_MODE_URL="https://modeanalytics.com/anki/reports/d8dd0c05c8a4"
env DAS_TABLE_NAME='{"OD": "das.odmessage", "COZMO": "das.cozmomessage"}'
env VECTOR_EVENTS="{'app'=>'https://modeanalytics.com/anki/reports/1c9b1ad383e9?run=now','robot'=>'https://modeanalytics.com/anki/reports/f874fe712488?run=now','cloud'=>'https://modeanalytics.com/anki/reports/ef268341fdda?run=now'}"

chdir /srv/repos/lita-slack-jira-hockeyapp-mode/
respawn
respawn limit 10 5

exec bash -c 'source /srv/.rvm/scripts/rvm; rvm ruby-2.1.0; cd /srv/repos/lita-slack-jira-hockeyapp-mode/workspace/bin; ./lita -c /srv/repos/lita-slack-jira-hockeyapp-mode/workspace/lita_config.rb' >> /tmp/upstart-hockeybot.log 2>&1
